# scripts-dev/com.exoplanet-hunter.refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- Weekly catalogue refresh (Saturday 09:00). Install:
       cp scripts-dev/com.exoplanet-hunter.refresh.plist ~/Library/LaunchAgents/
       launchctl load ~/Library/LaunchAgents/com.exoplanet-hunter.refresh.plist
     Optional: set NOTIFY_WEBHOOK_URL below for a Discord/Slack ping. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.exoplanet-hunter.refresh</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd /Users/ollie/Project/v2 &amp;&amp; /opt/anaconda3/envs/exoplanet-hunter-v2/bin/python orchestration/flows/refresh_pipeline.py --data-config full >> outputs/refresh-cron.log 2>&amp;1</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>6</integer>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <!--
  <key>EnvironmentVariables</key>
  <dict>
    <key>NOTIFY_WEBHOOK_URL</key>
    <string>https://discord.com/api/webhooks/...</string>
  </dict>
  -->
</dict>
</plist>
